# crontab
* * * * * /upload.sh /media/webcam.jpg ../../var/www/sobekkseter/src/public/img > /cron/cron.log 2>&1
