# crontab
* * * * * root /upload.sh test.txt . > ~/cron/cron.log 2>&1
